# scripts/templates/mdcz.service
[Unit]
Description=MDCz WebUI server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mdcz
Group=mdcz
# REPLACE_ME: install location of the portable bundle.
WorkingDirectory=/usr/lib/mdcz
Environment=PORT=3838
Environment=MDCZ_HOST=127.0.0.1
Environment=MDCZ_HOME=/var/lib/mdcz
Environment=MDCZ_WEB_DIST_DIR=/usr/lib/mdcz/web
ExecStart=/usr/bin/node /usr/lib/mdcz/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
